# startup.conf
[program:apache2]
command=/home/apache.sh
numprocs=1
autostart=true
autorestart=true

[program:mongod]
command=/usr/bin/mongod -f /etc/mongodb.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
autostart=true